# 99-qubes-tor-hook.rules
SUBSYSTEMS=="xen", KERNEL=="eth*", ACTION=="add", RUN+="/bin/systemctl --no-block qubes-tor.service"
